# racksview_app-1.1/etc/systemd/gstreamer-usb.service
[Unit]
Description=GStreamer capture from USB camera
After=network.target

[Service]
ExecStartPre=/usr/bin/rm -f /var/log/racksview/gstreamer-usb.log
ExecStart=/opt/racksview/start_gstreamer_usb.sh
Restart=on-failure
RestartSec=60
StandardOutput=null
StandardError=file:/var/log/racksview/gstreamer-usb.log
MemoryMax=256M

[Install]
WantedBy=multi-user.target
