# racksview_app-1.1/systemd/rvmanager.service
[Unit]
Description=Restart GStreamer capture services

[Service]
Type=oneshot
ExecStart=/opt/racksview/bin/cleanup.sh