# racksview_app-1.1/systemd/rvmanager.service
[Unit]
Description=Restart RacksView services and roll video storage

[Service]
Type=oneshot
ExecStart=/opt/racksview/scripts/cleanup.sh